# rootfs/etc/systemd/system/systemd-user-fix.service
[Unit]
Description=Systemd default user fix

[Service]
ExecStart=/usr/local/bin/systemd-user-fix.sh

[Install]
WantedBy=multi-user.target
